# conf/crontab
# Host crontab for the CloudVPS Compose Deploy toolkit, installed for the
# service account. The toolkit substitutes the ${...} placeholders via envsubst
# at provision time (the checkout path and branch are only known then).
#
# This carries the redeploy tick only. The app's own scheduled jobs (aggregate
# fills, linksearchtotal_collect, backups, ...) run inside the crons service off
# ./crontab, independent of how the stack is deployed.
SHELL=/bin/bash
PATH=${TOOLKIT_BIN}:${DEPLOY_ROOT}/bin:/usr/bin:/bin
DEPLOY_ROOT=${DEPLOY_ROOT}
COMPOSE_FILE=${DEPLOY_ROOT}/docker-compose.yml:${DEPLOY_ROOT}/docker-compose.deploy.yml

# git sync, pull, restart the stack on real change
*/5 * * * * deploy.sh ${ENV}
